# SOURCES/yana.cron
30 9,15 * * * yana /opt/yana/bin/run-nscan >/dev/null
